# crontab
# Nana Love Ads Bot
# 0 * * * * python3 /app/nana_ads_twitter.py
0 * * * * python3 /app/nana_ads_twitter_no_discord.py
